# cdrom.rules
#Detects new CDs and works with eject
KERNEL=="sr*", ACTION=="change", RUN+="${scriptPath}"
